# docs/systemd/shardline-api.service
[Unit]
Description=Shardline API role
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=shardline
Group=shardline
EnvironmentFile=/etc/shardline/shardline.env
Environment=SHARDLINE_SERVER_ROLE=api
ExecStart=/usr/local/bin/shardline serve --role api
ExecStartPre=/usr/local/bin/shardline config check
Restart=on-failure
RestartSec=5
StateDirectory=shardline
WorkingDirectory=/var/lib/shardline
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/shardline
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
